# systemd/job-radar.timer
[Unit]
Description=Run Local Job Radar four times daily

[Timer]
OnCalendar=*-*-* 09:00:00
OnCalendar=*-*-* 12:00:00
OnCalendar=*-*-* 16:00:00
OnCalendar=*-*-* 20:00:00
Persistent=true
Unit=job-radar.service

[Install]
WantedBy=timers.target
