# systemd/job-radar.service
[Unit]
Description=Local Job Radar Search Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/Desktop/JobSearch
ExecStart=%h/Desktop/JobSearch/.venv/bin/python %h/Desktop/JobSearch/main.py
StandardOutput=append:%h/Desktop/JobSearch/logs/systemd.log
StandardError=append:%h/Desktop/JobSearch/logs/systemd_error.log
